# 개발/사이니지파이/agent/systemd/memoria-signage.service
# 사이니지 에이전트 자동 실행 — 부팅 시 시작, 죽으면 즉시 재시작(검은화면 최소화).
#   설치: /opt/memoria 에 코드 배치 → 이 파일을 /etc/systemd/system/ 에 복사
#         sudo systemctl enable --now memoria-signage
[Unit]
Description=Memoria Signage Player Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/memoria
ExecStart=/usr/bin/python3 -m memoria_signage
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
